# federalism/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/log/supervisor/supervisord.pid
loglevel=info

[program:gunicorn]
# ✅ FIX: Added --worker-class gthread and --threads
command=gunicorn --bind 0.0.0.0:8000 --workers 2 --worker-class gthread --threads 4 --timeout 300 fed_api.wsgi:application
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/gunicorn.log
stderr_logfile=/var/log/supervisor/gunicorn_error.log

[program:celery]
# ✅ FIX: Added --pool=solo, --concurrency=1, and --max-tasks-per-child
command=celery -A fed_api worker --loglevel=info --pool=solo --concurrency=1 --max-tasks-per-child=50
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/celery.log
stderr_logfile=/var/log/supervisor/celery_error.log